# filiibot/rootfs/etc/services.d/filiibot/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Filii Lamberti Home Assistant Add-ons: Filiibot
# Starts the Filiibot
# ==============================================================================

# Set the current working dir to /opt
cd /opt/filiibot || bashio::exit.nok "Could not change directory to Filiibot"

bashio::log.info "Starting Filiibot..."

# And finally start the Filiibot
exec node bot.js
